# scripts/valheim-activity.service
[Unit]
Description="Valheim activity monitor"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/valheim-activity.sh
User=root
